# .docker/supervisor/all-in-one.conf
[supervisord]
nodaemon=true

[program:php]
priority=20
process_name=%(program_name)s
# command=php -d variables_order=EGPCS /var/www/html/artisan octane:start --server=swoole --host=0.0.0.0 --port=8000 --ansi
command=php-fpm -D
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:nginx]
priority=21
process_name=%(program_name)s
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:queue]
priority=21
process_name=%(program_name)s
command=php -d memory_limit=512M /var/www/html/artisan queue:work --queue=default,user,post,reward,crm,sms_product,sms_advertising --tries=3
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

# [program:cron]
# priority=22
# process_name=%(program_name)s
# command=/usr/sbin/crond -f -L 15
# autostart=true
# autorestart=true
# redirect_stderr=true
# stdout_logfile=/dev/stdout
# stderr_logfile=/dev/stderr
# stdout_logfile_maxbytes=0
# stderr_logfile_maxbytes=0
